# deploy/localweb-librespot.service
[Unit]
Description=librespot Spotify Connect receiver
After=network.target

[Service]
Type=simple
User=pi
ExecStartPre=/bin/bash -c 'test -p /tmp/librespot-pipe || mkfifo /tmp/librespot-pipe'
ExecStart=/usr/local/bin/librespot --name Drewtopia --backend pipe --device /tmp/librespot-pipe --bitrate 160 --volume-ctrl fixed
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
