# deploy/localweb.crontab
# localweb auto-deploy: check for updates every 5 minutes
*/5 * * * * /home/pi/localweb/deploy/update.sh
